# root/etc/systemd/system/docker.service.d/docker-systemd-override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd --dns 8.8.8.8 -H fd://
